# bin/instalador/scripts.arranque/tomcatgeopista6
#!/bin/sh
#
# tomcatgeopista6
#
# chkconfig: 2345 97 03
# description: 	Start up the Tomcat servlet engine.

export MAPSERVER_HOME=/usr/local/LocalGIS/mapserver
export LD_LIBRARY_PATH=$MAPSERVER_HOME/cgi-bin:$MAPSERVER_HOME/lib:$MAPSERVER_HOME/lib.ext:$LD_LIBRARY_PATH
export PROJ_LIB=$MAPSERVER_HOME/proj
 
RETVAL=$?
export JRE_HOME="/usr/java/jdk1.6.0_21/jre"
export CATALINA_HOME="/usr/local/apache-tomcat-6.0.33/"
export JAVA_HOME=/usr/java/jdk1.6.0_21
export CATALINA_OPTS="-Duser.dir=/tmp"

#Defino la ruta en la que se va a crear un fichero con el PID del proceso que se arranca
export CATALINA_PID=/usr/local/apache-tomcat-6.0.33/work/Catalina/catalina.pid

case "$1" in
 start)
        if [ -f $CATALINA_HOME/bin/startup.sh ];
          then
	    echo $"Starting Tomcat"
            /bin/su tomcat $CATALINA_HOME/bin/startup.sh
            echo "    ...done."                    
        fi
	;;
  stop)
        if [ -f $CATALINA_HOME/bin/shutdown.sh ]; then
	    echo "  * Stopping Tomcat"
	    if [ "$2" = "-force" ]; then
                # Se ejecuta el shutdown - force con el usuario "tomcat"
                /bin/su -c "sh /usr/local/apache-tomcat-6.0.33/bin/shutdown.sh $2" tomcat
            else
                # Se ejecuta un shutdown sin force con el usuario "tomcat"
                /bin/su tomcat $CATALINA_HOME/bin/shutdown.sh
            fi   
            echo "    ...done."
        fi
 	;;
  *)
 	echo $"Usage: $0 {start | stop [-force]}"
	exit 1
	;;
esac

exit $RETVAL
